# recipes-bsp/telit/files/modem-at.service
[Unit]
Description=Modem AT

BindsTo=dev-modem_at.device
After=dev-modem_at.device

[Service]
ExecStart=/usr/bin/runmodem.sh

[Install]
WantedBy=multi-user.target
